# etc/init.d/bot-telegram-area31
#!/sbin/openrc-run
#set -x

USER="morfetico"
GROUP="alta-linguagem"
BOTDIR="$(getent passwd ${USER} | cut -d: -f6)/${GROUP}"

depend() {
need net
}

start() {
ebegin "Iniciando o Bot-Telegram \"${USER}\"/${GROUP}"
echo -e "while true; do\ncd ${BOTDIR} ; python3.12 ${BOTDIR}/bot-telegram.py\ndone" > ${BOTDIR}/start.sh
chmod +x ${BOTDIR}/start.sh
su ${USER} -c "nohup ${BOTDIR}/start.sh >> ${BOTDIR}/bot-telegram.log & &> /dev/null"  &> /dev/null
echo -e "" &> /dev/null
}

stop() {
ebegin "Parando o Bot-Telegram \"${USER}\""
pid_info=$(ps aux | grep ${GROUP}|grep "start.sh"|grep -v grep)
pids=$(echo "$pid_info" | awk '{print $2}' | head -n 1)
if [ -z "$pids" ]; then
  echo "Process not found."
  #exit 1
  echo
fi
for pid in $pids; do
kill $pid
done
echo

pid_info=$(ps aux | grep ${GROUP}|grep "bot-telegram.py"|grep -v grep)
pids=$(echo "$pid_info" | awk '{print $2}' | head -n 1)
if [ -z "$pids" ]; then
  echo "Process not found."
 # exit 1
  echo
fi
for pid in $pids; do
kill $pid
done
echo -e "" &> /dev/null
}

status() {
ebegin "Processo em execução"

echo -e "\nUSER: ${USER}"

CHECK=$(pgrep -f "bash ${BOTDIR}/start.sh"|wc -l)
if [ "$CHECK" -ne 0 ]; then
	echo -e "start.sh PID: $(pgrep -f "bash ${BOTDIR}/start.sh")"
else
	echo -e "start.sh NOT running. Please restart bot-telegram"
fi

echo -e ""
}
