# etc/init.d/bot-telegram-area31
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2
#set -x

USER="morfetico"
GROUP="area31"
BOTDIR="$(getent passwd ${USER} | cut -d: -f6)/${GROUP}"

depend() {
need net
}

start() {
ebegin "Iniciando o Bot-Telegram \"${USER}\"/${GROUP}"
echo -e "while true; do\ncd ${BOTDIR} ; python ${BOTDIR}/bot-telegram.py\ndone" > ${BOTDIR}/start.sh
chmod +x ${BOTDIR}/start.sh
mv ${BOTDIR}/bot-telegram.log2 ${BOTDIR}/bot-telegram.log3 2> /dev/null
mv ${BOTDIR}/bot-telegram.log ${BOTDIR}/bot-telegram.log2 2> /dev/null
su ${USER} -c "nohup ${BOTDIR}/start.sh > ${BOTDIR}/bot-telegram.log & &> /dev/null"  &> /dev/null
echo -e "" &> /dev/null
}

stop() {
ebegin "Parando o Bot-Telegram \"${USER}\""
pid_info=$(ps aux | grep ${GROUP}|grep "start.sh"|grep -v grep)
pids=$(echo "$pid_info" | awk '{print $2}' | head -n 1)
if [ -z "$pids" ]; then
  echo "Process not found."
  exit 1
fi
for pid in $pids; do
kill $pid
done

pid_info=$(ps aux | grep ${GROUP}|grep "bot-telegram.py"|grep -v grep)
pids=$(echo "$pid_info" | awk '{print $2}' | head -n 1)
if [ -z "$pids" ]; then
  echo "Process not found."
  exit 1
fi
for pid in $pids; do
kill $pid
done
echo -e "" &> /dev/null
}

status() {
ebegin "Processo em execução"

echo -e "\nUSER: ${USER}"

CHECK=$(pgrep -f "bash ${BOTDIR}/start.sh"|wc -l)
if [ "$CHECK" -ne 0 ]; then
	echo -e "start.sh PID: $(pgrep -f "bash ${BOTDIR}/start.sh")"
else
	echo -e "start.sh NOT running. Please restart bot-telegram"
fi

#CHECK2=$(pgrep -u${USER} telegram-cli|wc -l)
#if [ "$CHECK2" -ne 0 ]; then
#	echo -e "telegram-cli PID: $(pgrep -u${USER} telegram-cli)"
#else
#	echo -e "telegram-cli NOT running. Please restart bot-telegram"
#fi
echo -e ""
}
